# config/services/certbot-renew.timer
[Unit]
Description=certbot renew timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=cloudbox.target
